# autostart/lxqt-globalkeyshortcuts.desktop.in
[Desktop Entry]
Type=Application
TryExec=lxqt-globalkeysd
Exec=lxqt-globalkeysd
OnlyShowIn=LXQt;
X-LXQt-Module=true
X-LXQt-X11-Only=true

#TRANSLATIONS_DIR=translations
